# Institute/Sisop/Praktikum/1/task-3/crontabs
*/2 * * * * /home/user/praktikum-modul-1-c11-2/task-3/automation.sh
*/10 * * * * /home/user/praktikum-modul-1-c11-2/task-3/download.sh
0 */2 * * * /home/user/praktikum-modul-1-c11-2/task-3/archive.sh